# home/dot_local/share/portable-laptop-power/root/thinkpad-power-source-policy.service
[Unit]
Description=Apply laptop power policy for the current AC state
After=power-profiles-daemon.service iwd.service
Wants=power-profiles-daemon.service
ConditionPathExists=/sys/class/power_supply/AC/online

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/thinkpad-power-source-policy

[Install]
WantedBy=multi-user.target
