# home/dot_config/systemd/user/decent-angl-restic-forget-prune.service.tmpl
[Unit]
Description=Decent Angl restic forget and prune
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/restic-scheduled-job forget-prune
