# home/dot_config/systemd/user/decent-angl-restic-restore-smoke.service.tmpl
[Unit]
Description=Decent Angl restic restore smoke test
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/restic-scheduled-job restore-smoke
