# home/dot_config/systemd/user/decent-angl-restic-check.service.tmpl
[Unit]
Description=Decent Angl restic integrity check
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/restic-scheduled-job check
