# home/dot_config/systemd/user/decent-angl-restic-backup.service.tmpl
[Unit]
Description=Decent Angl restic backup
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/restic-scheduled-job backup
