# home/dot_config/systemd/user/omarchy-battery-monitor.service
[Unit]
Description=Omarchy Battery Monitor Check
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/share/omarchy/bin/omarchy-battery-monitor
Environment=PATH=%h/.local/share/omarchy/bin:%h/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/bin:/usr/sbin:/sbin:/bin
Environment=DISPLAY=:0
LogLevelMax=warning
